# home/.config/systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata

[Service]
Type=simple
ExecStart=/usr/bin/kanata --cfg %h/.kanata.conf
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target